# home/.config/systemd/dconf-icon-monitor.path
[Unit]
Description=Theo doi bien dong thu muc Desktop de ghi lai toa do icon

[Path]
# Theo dõi mọi hành động chỉnh sửa, thêm, xóa file trong thư mục Desktop
PathChanged=%h/Desktop

[Install]
WantedBy=gnome-session.target
