# home/.config/systemd/flatpak-upgrade.timer
[Unit]
Description=Flatpak Automatic Update Timer

[Timer]
OnUnitActiveSec=4h
OnBootSec=15m

[Install]
WantedBy=timers.target
