# home/.config/systemd/flatpak-upgrade.service
[Unit]
Description=Flatpak Automatic Update

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update -y
