# home/.config/systemd/wine-theme-sync.service
[Unit]
Description=Wine Auto Theme Sync

[Service]
ExecStart=/opt/wine-themes/wine-auto-sync
Restart=always

[Install]
WantedBy=default.target
